# poll.service
[Unit]
Description=Constantly checks for files in a directory
